# scripts/nick2-gate-tunnel.service
[Unit]
Description=Cloudflared HTTPS tunnel for Nick2 gate-chat API
After=network.target nick2-gate-chat.service

[Service]
Type=simple
ExecStart=/tmp/cloudflared tunnel --url http://127.0.0.1:8788
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target